# provisioning/backup/files/systemd/backup.service
[Unit]
Description=Backup www directory

[Service]
Type=oneshot
ExecStart=/opt/backup.sh
