# provisioning/backup/files/systemd/backup-sql.service
[Unit]
Description=Backup sql directory on glusterfs mount-point


[Service]
Type=oneshot
ExecStart=/opt/backup-sql.sh
