# provisioning/backup/files/systemd/backup-elk.service
[Unit]
Description=Backup elk directory on glusterfs mount-point


[Service]
Type=oneshot
ExecStart=/opt/backup-elk.sh
