# provisioning/backup/files/mysql-replica/systemd/backup-sql.service
[Unit]
Description=SQL dump MariaDB

[Service]
Type=oneshot
ExecStart=/opt/backup-sql.sh
